# llll.conf
[program:llll]
command=python3 /home/curator2/llll/main.py
autostart=true
autorestart=true
stderr_logfile=/home/curator2/llll/test.err.log
stdout_logfile=/home/curator2/llll/test.out.log